# deploy/intellexia-mcp.service
[Unit]
Description=IntellexIA MCP Server (OAuth)
After=network.target

[Service]
User=root
WorkingDirectory=/sites/intellexia
ExecStart=/root/.local/bin/uv run python mcp_server/server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8001
# O domínio NÃO entra aqui: vem do .env da instalação (APP_PUBLIC_URL), que o
# main.py carrega por cima do ambiente do processo — um valor nesta unit seria
# sobrescrito pelo .env e só geraria confusão sobre quem manda.

[Install]
WantedBy=multi-user.target
